# bambuddy/rootfs/etc/services.d/bambuddy/run
https://github.com/Spegeli/homeassistant-app-bambuddy/tree/main/bambuddy/rootfs/etc
